# backend/deploy/systemd/mingxintai-reminders.service
[Unit]
Description=Mingxintai reminder scan
After=network.target mingxintai.service

[Service]
Type=oneshot
EnvironmentFile=/etc/mingxintai/backend.env
ExecStart=/bin/sh -lc 'cfg="$(mktemp /run/mingxintai-reminders.XXXXXX)"; trap "rm -f \"$cfg\"" EXIT; chmod 600 "$cfg"; printf "%s\n" "request = POST" "url = http://127.0.0.1:8000/api/reminders/scan" "header = Content-Type: application/json" "header = X-Mingxin-Token: ${MINGXIN_API_TOKEN}" "data = {\"soon_hours\":168,\"dry_run\":false}" "fail" "silent" "show-error" > "$cfg"; curl --config "$cfg" >/var/log/mingxintai-reminders.log'
